# config/systemd/modulehub-cms.service
[Unit]
Description=ModuleHub CMS — Modular Docker CMS
After=network.target docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/modulehub-cms
EnvironmentFile=/opt/modulehub-cms/.env
ExecStart=/usr/bin/node /opt/modulehub-cms/dist/core/src/server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
